# system/templates/btrfs-scrub.service.j2
# Managed by infra repo
[Unit]
Description=btrfs scrubber

[Service]
Type=oneshot
ExecStart=btrfs scrub start /
