# system/templates/paccache.timer.j2
# Managed by infra repo
[Unit]
Description=pacman cache cleaner

[Timer]
OnCalendar=Sat, 05:00:00

[Install]
WantedBy=timers.target
