# system/templates/sysctl-boot.service.j2
# Managed by infra repo
[Unit]
Description=sysctl apply

[Service]
Type=oneshot
ExecStart=sysctl -p
