# system/templates/sysctl-boot.timer.j2
# Managed by infra repo
[Unit]
Description=sysctl apply on boot

[Timer]
OnBootSec=1s

[Install]
WantedBy=timers.target
